# wall.socket
[Unit]
Description=NodeWall for LinBPQ

[Socket]
ListenStream=6001
Accept=yes

[Install]
WantedBy=sockets.target
